# systemd/hibp-checker.timer
[Unit]
Description=HIBP Breach Checker - Scheduled Monitoring
Documentation=https://github.com/greogory/hibp-checker
Requires=hibp-checker.service

[Timer]
# Run daily at 2 AM
OnCalendar=*-*-* 02:00:00

# Run 15 minutes after boot if the 2 AM schedule was missed
OnBootSec=15min

# Run on boot if missed (e.g., computer was off)
Persistent=yes

# Randomize start time by up to 30 minutes to avoid API rate limit conflicts
RandomizedDelaySec=30min

# Ensure timer is accurate
AccuracySec=1min

[Install]
WantedBy=timers.target
