# dashboard/systemd/hibp-dashboard.service
[Unit]
Description=HIBP Dashboard Web Server
Documentation=https://github.com/greogory/hibp-checker
After=network.target

[Service]
Type=simple

# IMPORTANT: Update this path to your hibp-project/dashboard installation directory
# Example: WorkingDirectory=%h/hibp-project/dashboard
WorkingDirectory=HIBP_PROJECT_DIR/dashboard

# IMPORTANT: Update this path to match WorkingDirectory
ExecStart=HIBP_PROJECT_DIR/dashboard/start-dashboard.sh

Restart=on-failure
RestartSec=5

# Output (XDG-compliant location)
StandardOutput=append:%h/.local/share/hibp-checker/dashboard.log
StandardError=append:%h/.local/share/hibp-checker/dashboard.error.log

# Security hardening (match hibp-checker.service)
ProtectSystem=strict
ProtectHome=read-only
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=default.target
